# archiso/airootfs/etc/profile.d/frostbite-gaming.sh
export MANGOHUD_CONFIGFILE="${MANGOHUD_CONFIGFILE:-$HOME/.config/MangoHud/MangoHud.conf}"
export STEAM_EXTRA_COMPAT_TOOLS_PATHS="${STEAM_EXTRA_COMPAT_TOOLS_PATHS:-$HOME/.steam/root/compatibilitytools.d}"
